# conf/upstart.conf
description "CodeCube Server"
author      "Harry Marr"

start on startup
stop on shutdown

script
    exec /home/harry/runner/runner
end script

